# dist/macos/payload/Library/LaunchDaemons/com.graylog.collector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.graylog.collector</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/graylog-collector/bin/graylog-collector</string>
      <string>supervisor</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/var/log/graylog-collector.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/graylog-collector.err.log</string>
  </dict>
</plist>
